# mailpod-entrypoint.service
[Unit]
Description=Mailpod entrypoint

[Service]
Type=simple
ExecStart=/usr/local/bin/mailpod-entrypoint

[Install]
WantedBy=default.target
